# SOURCES/finger@.service
[Unit]
Description=Finger Per-Connection Server

[Service]
User=nobody
ExecStart=-/usr/sbin/efingerd
StandardInput=socket
