# SOURCES/nncp.socket
[Unit]
Description=NNCP Daemon Activation Socket
Documentation=info:nncp

[Socket]
ListenStream=0.0.0.0:1199

[Install]
WantedBy=sockets.target
