# SOURCES/nncp.service
[Unit]
Description=NNCP Daemon
Documentation=info(nncp)

[Service]
ExecStart=/usr/bin/nncp-daemon
StandardInput=socket

[Install]
Also=nncp.socket
